# php/supervisord-worker.ini
[program:queue-work]
command=php artisan queue:work
directory=/data/www/mansion
autostart=true
autorestart=true
redirect_stderr=true
startretries=1
user=www-data
numprocs=5
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
